# batnotifyd.service
[Unit]
Description=batnotifyd

[Service]
ExecStart=/usr/local/bin/batnotifyd
Restart=always

[Install]
WantedBy=default.target
